# deploy/unifi-door-watcher.service
[Unit]
Description=UniFi Access door watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=unifi-door-watcher
EnvironmentFile=/etc/unifi-door-watcher/env
ExecStart=/opt/unifi-door-watcher/.venv/bin/python -m unifi_door_watcher --config /etc/unifi-door-watcher/config.toml
Restart=always
RestartSec=5
ProtectSystem=strict
# `read-only` (not `true`) so the service can read its own code when it's
# deployed under /home/forge/… by Laravel Forge. It still can't write to
# any home directory — that's the substantive protection we want.
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=

[Install]
WantedBy=multi-user.target
